# systemd/user/stenogit@.service.in
[Unit]
Description=Stenogit commit for %i
ConditionPathExists=%h/.config/stenogit/%i.conf

[Service]
Type=oneshot
EnvironmentFile=%h/.config/stenogit/%i.conf
Environment=INSTANCE=%i
ExecStart=@BINDIR@/stenogit-commit
